# schlage-app.service
[Unit]
Description=Schlage Lock Management App
After=network.target

[Service]
Type=simple
User=schlage
Group=schlage
WorkingDirectory=/opt/schlage-app
Environment="PATH=/opt/schlage-app/venv/bin"
ExecStart=/opt/schlage-app/start_server.py
Restart=on-failure
RestartSec=5s
StandardOutput=append:/opt/schlage-app/logs/app.log
StandardError=append:/opt/schlage-app/logs/app.log

[Install]
WantedBy=multi-user.target
